# docker/config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
username=nobody
password=nobody

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root                    ; (default is current user, required if root)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:php-fpm7]
command=/usr/sbin/php-fpm7.4 --nodaemonize --fpm-config=/etc/php/7.4/fpm/pool.d/www.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler

[program:writerzen-worker]
command=php /var/www/html/artisan queue:work --queue=high,medium,default
process_name=%(program_name)s_%(process_num)02d
numprocs=20 
priority=999
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log

[program:writerzen-broadcast]
command=php /var/www/html/artisan queue:work --queue=broadcast --max-time=36000 --timeout=36000
process_name=%(program_name)s_%(process_num)02d
numprocs=2
priority=998
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log

[program:writerzen-cluster]
command=php /var/www/html/artisan queue:work --queue=cluster --max-time=36000 --timeout=36000
process_name=%(program_name)s_%(process_num)02d
numprocs=2 
priority=997
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log

[program:writerzen-cluster-large]
command=php /var/www/html/artisan queue:work --queue=large_cluster --max-time=36000 --timeout=36000
process_name=%(program_name)s_%(process_num)02d
numprocs=1
priority=997 
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log

[program:writerzen-balancer-1]
command=php /var/www/html/artisan queue:work --queue=balancer_1,balancer_2,balancer_3 --max-time=60 --timeout=60
process_name=%(program_name)s_%(process_num)02d
numprocs=4
priority=996 
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log

[program:writerzen-balancer-2]
command=php /var/www/html/artisan queue:work --queue=balancer_2,balancer_1,balancer_3 --max-time=60 --timeout=60
process_name=%(program_name)s_%(process_num)02d
numprocs=4 
priority=996 
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log

[program:writerzen-balancer-3]
command=php /var/www/html/artisan queue:work --queue=balancer_3,balancer_2,balancer_1 --max-time=60 --timeout=60
process_name=%(program_name)s_%(process_num)02d
numprocs=4
priority=996 
autostart=true
autorestart=true  
startsecs=1
startretries=3000
;user=apache
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log